# tele-aria2/tele-aria2.service
[Unit]
Description=File Browser
After=network.target

[Service]
ExecStart=tele-aria2 -c /root/.config.json

[Install]
WantedBy=multi-user.target
